# .github/systemd/app-tl-ai.service
[Unit]
Description=true-love-ai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# 首次部署（只需执行一次, 请手动执行）：
# sudo install -m 0644 app-tl-ai.service /etc/systemd/system/app-tl-ai.service
# sudo systemctl daemon-reload
# sudo systemctl enable app-tl-ai
# sudo systemctl start app-tl-ai

# === pm2: cwd ===
WorkingDirectory=/root/true-love-n/true-love-ai

# === pm2: env ===
Environment=APP_ENV=prod

# === pm2: interpreter + script ===
ExecStart=/root/true-love-n/true-love-ai/.venv/bin/python main.py

# === pm2: autorestart ===
Restart=always
RestartSec=2

# === pm2: time + log_date_format ===
# journald 自带时间戳（ISO-8601），比 pm2 稳
StandardOutput=journal
StandardError=journal

# 建议指定运行用户
User=root
Group=root

# 优雅退出
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target